# config/services/vault-agent/systemd/abhaile-secrets-ready.service
[Unit]
Description=Abhaile secrets/templates are rendered (global gate)

[Service]
Type=oneshot
ExecStart=/bin/true
RemainAfterExit=yes
