# config/services/coredns-omada/systemd/coredns-omada-env.service
[Unit]
Description=Restart CoreDNS after Omada env update
[Service]
Type=oneshot
ExecCondition=/bin/sh -c 'test ! -e /etc/coredns/omada.env || ! cmp -s /srv/vault/agent/out/coredns-omada.env /etc/coredns/omada.env'
ExecStartPre=/usr/bin/install -d -m 0750 -o root -g coredns /etc/coredns
ExecStart=/usr/bin/install -m 0640 -o root -g coredns /srv/vault/agent/out/coredns-omada.env /etc/coredns/omada.env
ExecStartPost=/bin/systemctl --no-block restart coredns.service
