# out_bak/rendered/phobos/services/coredns-filtered/etc/systemd/system/coredns-zones.service
[Unit]
Description=Reload CoreDNS configuration
Documentation=https://coredns.io
# Only reload if CoreDNS is actually running
PartOf=coredns.service

[Service]
Type=oneshot
# Graceful reload via SIGUSR1 (CoreDNS reloads config without dropping connections)
ExecStart=/usr/bin/systemctl kill --signal=SIGUSR1 coredns.service
# Verify reload succeeded
ExecStartPost=/usr/bin/sleep 1
ExecStartPost=/usr/bin/systemctl is-active coredns.service

[Install]
WantedBy=multi-user.target
