# out_bak/rendered/phobos/services/coredns-filtered/etc/systemd/system/coredns-omada-env.path
[Unit]
Description=Watch CoreDNS Omada env for changes
[Path]
PathChanged=/srv/vault/agent/out/coredns-omada.env
[Install]
WantedBy=multi-user.target
