# systemd/abhaile-gitops-apply@.service
[Unit]
Description=Abhaile GitOps apply for %i (root)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/opt/abhaile
ExecStart=/opt/abhaile/tools/apply/apply.sh --skip-render --apply %i
ExecStartPost=/usr/bin/rm -f /var/lib/abhaile/rendered/.apply_ready

[Install]
WantedBy=multi-user.target
