# systemd/abhaile-gitops@.timer
[Unit]
Description=Periodic GitOps apply for %i
Conflicts=abhaile-gitops@%i.service

[Timer]
OnBootSec=10m
OnUnitActiveSec=15m
RandomizedDelaySec=120
Persistent=true
# L4: Prevent stacked processes - do not queue runs while service is active
Unit=abhaile-gitops@%i.service

[Install]
WantedBy=timers.target
