# config/services/authelia/systemd/authelia-config.path
[Unit]
Description=Watch Authelia configuration (rendered) for changes

[Path]
PathChanged=/srv/vault/agent/out/authelia.configuration.yml

[Install]
WantedBy=multi-user.target
