# config/hosts/common/systemd/abhaile-runner.service
[Unit]
Description=Abhaile GitOps Runner
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/abhaile/scripts/abhaile-runner
User=abhaile
Group=abhaile
WorkingDirectory=/opt/abhaile
Environment=PATH=/opt/abhaile/.venv/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=-/etc/abhaile/runner.env
StandardOutput=journal
StandardError=journal
