# out_bak/rendered/phobos/services/omada-controller/etc/systemd/system/rebuild-omada-cert.service
[Unit]
Description=Rebuild Omada TLS full chain and restart Omada Controller
After=caddy-internal.service

[Service]
Type=oneshot
ExecStart=/opt/abhaile/tools/bash/rebuild-omada-cert.sh
