# manual_config/phobos/etc/containers/systemd/caddy-dmz.container
[Unit]
Description=Caddy (DMZ)
After=network-online.target
Wants=network-online.target

[Container]
DNS=172.20.20.236
EnvironmentFile=/etc/caddy/dns.env
Image=caddy-dmz.build
Network=dmz.network
PodmanArgs=--ip=172.20.100.200
Timezone=local
Volume=/srv/caddy/dmz/Caddyfile:/etc/caddy/Caddyfile:Z
Volume=caddy-dmz-data.volume:/data:Z
Volume=caddy-dmz-config:/config:Z
Volume=host-certs.volume:/etc/ssl/certs:ro

[Service]
Restart=on-failure
RestartSec=3s

[Install]
WantedBy=multi-user.target
