# manual_config/deimos/etc/systemd/system/coredns-omada-env.service
[Unit]
Description=Restart CoreDNS after Omada env update
[Service]
Type=oneshot
ExecStart=/bin/systemctl restart coredns.service
