# manual_config/deimos/etc/systemd/system/coredns-omada-install.service
[Unit]
Description=Install CoreDNS (Omada build) to /usr/local/bin
After=network-online.target coredns-omada-build.service
Requires=coredns-omada-build.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'podman run --rm --pull=never --entrypoint=/bin/sh -v /usr/local/bin:/out:Z localhost/coredns-omada:latest -c "cp /coredns /out/coredns && chmod 0755 /out/coredns"'

[Install]
WantedBy=multi-user.target
