# manual_config/phobos/etc/systemd/system/authelia-redis-conf.service
[Unit]
Description=Restart redis after config change
[Service]
Type=oneshot
ExecStart=systemctl restart authelia-app-redis.service
