# manual_config/phobos/etc/systemd/system/authelia-redis-conf.path
[Unit]
Description=Watch Redis config from Vault Agent

[Path]
PathChanged=/srv/vault/agent/out/authelia-redis.conf

[Install]
WantedBy=multi-user.target
