# manual_config/phobos/etc/systemd/system/authelia-config.service
[Unit]
Description=Restart Authelia after configuration change
[Service]
Type=oneshot
ExecStart=systemctl restart authelia-app-authelia.service
