# manual_config/phobos/etc/systemd/system/caddy-dns-desec.service
[Unit]
Description=Restart Caddy (DMZ) after token changes
[Service]
Type=oneshot
ExecStart=systemctl restart caddy-dmz.service
