# manual_config/phobos/etc/systemd/system/vault-ready.service
[Unit]
Description=Vault Ready (unsealed)
Requires=vault.service abhaile-vault-unseal.service
After=vault.service abhaile-vault-unseal.service
