# config/services/omada-controller/systemd/omada-controller-env.service
[Unit]
Description=Copy Omada Controller environment from Vault Agent
After=abhaile-secrets-ready.service
Requires=abhaile-secrets-ready.service

[Service]
Type=oneshot
ExecCondition=/bin/sh -c 'test ! -e /etc/omada-controller/omada-controller.env || ! cmp -s /srv/vault/agent/out/omada-controller.env /etc/omada-controller/omada-controller.env'
ExecStartPre=/usr/bin/install -d -m 0750 -o root -g root /etc/omada-controller
ExecStart=/usr/bin/install -m 0640 -o root -g root /srv/vault/agent/out/omada-controller.env /etc/omada-controller/omada-controller.env
ExecStartPost=/usr/bin/systemctl --no-block try-restart omada-controller-app-omada-controller.service
