# config/services/omada-controller/systemd/omada-controller-env.path
[Unit]
Description=Watch Omada Controller environment from Vault Agent

[Path]
PathChanged=/srv/vault/agent/out/omada-controller.env

[Install]
WantedBy=multi-user.target
