# config/services/omada-controller/systemd/rebuild-omada-cert.service
[Unit]
Description=Rebuild Omada TLS full chain and restart Omada Controller
After=caddy-internal.service omada-controller-env.service
Requires=omada-controller-env.service

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/local/lib/abhaile/tools/rebuild-omada-cert.sh
