# config/services/omada-controller/systemd/rebuild-omada-cert.service
[Unit]
Description=Rebuild Omada TLS full chain and restart Omada Controller
After=caddy-internal.service

[Service]
Type=oneshot
ExecStart=/usr/local/lib/abhaile/tools/rebuild-omada-cert.sh
