# out_bak/rendered/phobos/services/vault-agent/etc/systemd/system/abhaile-secrets-ready.path
[Unit]
Description=Watch Vault Agent sentinel for secrets readiness

[Path]
PathExists=/srv/vault/agent/out/.ready

[Install]
WantedBy=multi-user.target
